# attendance-tracker.service
[Unit]
Description=Attendance Tracker Service
After=network.target

[Service]
Type=simple
User=ivo
WorkingDirectory=/home/ivo/.openclaw/workspace/attendance-tracker
ExecStart=/home/ivo/.openclaw/workspace/attendance-tracker/.venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
